# backend.service
[Unit]
Description=Online Chess Game Backend Service
Documentation=https://github.com/yourusername/online-chess-game
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/online-chess-game
ExecStart=/usr/bin/node /opt/online-chess-game/server.cjs
Restart=always
RestartSec=10
StandardOutput=append:/opt/online-chess-game/logs/backend.log
StandardError=append:/opt/online-chess-game/logs/backend-error.log
Environment=NODE_ENV=production
Environment=PORT=3001

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/online-chess-game/logs

[Install]
WantedBy=multi-user.target
